# ops/systemd/tj-main.service
[Unit]
Description=TractionJackson Main Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/traction-jackson
EnvironmentFile=/opt/traction-jackson/.env
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tj-main

[Install]
WantedBy=multi-user.target
